# ml/scripts/com.projectgrader.psa-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for the daily data-flywheel cycle.

  Sequence (per ml/scripts/daily_cycle.sh):
    1. psa_daily_ingest         — drain the PSA Public API daily budget.
    2. embed_references         — embed any new manufacturer references.
    3. accumulate_psa_exemplars — append today's matched PSA records as
                                  per-variant authentic exemplars to
                                  reference_embeddings.npz, bumping
                                  embedding-anomaly's confidence ramp.

  Label is unchanged from the original ingest-only job
  (`com.projectgrader.psa-ingest`) so existing installations keep working
  after a `launchctl unload && load` cycle — no need to rename files in
  ~/Library/LaunchAgents/.

  Install:
    cp ml/scripts/com.projectgrader.psa-ingest.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.projectgrader.psa-ingest.plist

  Trigger immediately (test):
    launchctl start com.projectgrader.psa-ingest

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.projectgrader.psa-ingest.plist
    rm ~/Library/LaunchAgents/com.projectgrader.psa-ingest.plist

  Edit the hardcoded PROJECT_ROOT below if you fork or relocate the project.
  See ml/scripts/README.md for full setup including the .env file.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.projectgrader.psa-ingest</string>

    <!-- The wrapper handles cd into the project root (default
         /Users/danielbelmeur/Github/project-grader; override via
         PROJECT_ROOT env) and sources .env itself, so the plist just
         hands off. Edit the absolute path here AND the default in
         daily_cycle.sh if you fork or relocate the project. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/danielbelmeur/Github/project-grader/ml/scripts/daily_cycle.sh</string>
    </array>

    <!-- 09:00 local time daily. UTC midnight resets the PSA daily budget,
         so any time after that works; 09:00 local gives you observable
         output during the day. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't fire on first load — the first run should be deliberate
         (`launchctl start ...`) so we can watch the logs. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/danielbelmeur/psa_data/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/danielbelmeur/psa_data/launchd.stderr.log</string>
</dict>
</plist>
